# conteiners/Dockerfiles/centos-supervisiord-nginx-web2py/cloudportal/supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=debug               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[program:gunicorn]
command=/bin/gunicorn wsgihandler:application -c /var/www/web2py/gunicorn.conf.py
directory=/var/www/web2py ;
user=root ;
autostart=true ;
autorestart=true ;
redirect_stderr=True ;
environment=PYTHONPATH="/bin/python2.7",MYSQL_CONNECTION="" ;
stdout_logfile = /var/log/supervisor/gunicorn.log ; 
stderr_logfile = /var/log/supervisor/gunicorn_err.log ;

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
